# docker/templates/php7-nginx-mysql57/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/sshd.err.log
stdout_logfile=/var/log/supervisor/sshd.out.log