# docker/templates/php7.1-nginx-mysql5.7/supervisor/conf.d/php.conf
[program:php]
command=php-fpm
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/php.err.log
stdout_logfile=/var/log/supervisor/php.out.log